# supervisord.conf
[supervisord]
nodaemon=true

[program:hn_ocr]
autorestart=True
autostart=True
process_name=hn_ocr
command=python ./hn_ocr/server.py
directory=/
user=root
